# stow-files/webserver/etc/systemd/system/ustreamer.service
[Unit]
Description=Lightweight HTTP MJPG streamer for webcam
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
ExecStart=/usr/bin/ustreamer --device=/dev/video0 --host=127.0.0.1 --port=8081 --format=YUYV --workers=2 --persistent
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
