# windows-11/dotfiles/powershell/profile.d/03-extra.ps1
# Entorno Go (equivalente a .bashrc.d/05-extra.sh).
$env:GOPATH = "$HOME\dev\go"
if ($env:PATH -notlike "*$env:GOPATH\bin*") {
    $env:PATH = "$env:GOPATH\bin;$env:PATH"
}
