# windows-11/dotfiles/powershell/profile.d/04-work-aliases.ps1
# Aliases de trabajo (placeholder).
# Los de Fedora (vm/vup/vhalt/vmmount) son especificos de Vagrant/sshfs en Linux y
# viven en WSL, no en PowerShell nativo. Agregar aca los equivalentes Windows si hacen falta.
